# import/90-import-photos.rules
ACTION=="add", SUBSYSTEM=="block", RUN+="/usr/local/sbin/photobox_import.sh"
